# ops/systemd/madcore-certbot-renew.timer
[Unit]
Description=Daily MADCORE Let's Encrypt renewal check

[Timer]
OnCalendar=*-*-* 04:15:00
RandomizedDelaySec=30m
Persistent=true
Unit=madcore-certbot-renew.service

[Install]
WantedBy=timers.target
